# ghostbsd-settings/repos/current/ports/common/live/pkg-plist
share/desktopbsd/common-live-settings-desktopbsd/base/backup/etc/rc.d/mountcritlocal
share/desktopbsd/common-live-settings-desktopbsd/base/override/etc/rc.d/mountcritlocal
share/desktopbsd/common-live-settings-desktopbsd/etc/rc.d/unionfs
share/desktopbsd/common-live-settings-desktopbsd/config-live-settings
share/desktopbsd/common-live-settings-desktopbsd/etc/rc.d/uzip
